# config/webserver/systemd/backupdb.timer
[Unit]
Description=Back up sqlite3 database daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target